# root/etc/s6-overlay/s6-rc.d/certbot-prepare/up
exec s6-setuidgid mrmeeb /certbot-prepare.sh